# inhatc/deploy/inhatc.service
[Unit]
Description=Inhatc Spring Boot Application
After=network.target mysql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=DEPLOY_PATH/app
ExecStart=/usr/bin/java -jar -Dspring.profiles.active=prod DEPLOY_PATH/app/inhatc-0.0.1-SNAPSHOT.jar
ExecStop=/bin/kill -15 $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inhatc

# 환경 변수 설정
Environment="JAVA_HOME=/usr/lib/jvm/default-java"

# 데이터베이스 연결 정보 (보안상 필수)
Environment="DB_URL=jdbc:mysql://localhost:3306/member?useSSL=false&allowPublicKeyRetrieval=true&serverTimezone=Asia/Seoul"
Environment="DB_USERNAME=your_database_username"
Environment="DB_PASSWORD=your_secure_database_password"

# 파일 업로드 경로
Environment="UPLOAD_BASE_DIR=DEPLOY_PATH/uploads"
Environment="UPLOAD_POSTS_DIR=DEPLOY_PATH/uploads/posts"
Environment="UPLOAD_PROFILE_DIR=DEPLOY_PATH/uploads/profiles"
Environment="UPLOAD_MESSAGES_DIR=DEPLOY_PATH/uploads/messages"

# 메모리 설정 (필요에 따라 조정)
Environment="JAVA_OPTS=-Xms512m -Xmx1024m -XX:+UseG1GC"

[Install]
WantedBy=multi-user.target



